# root/etc/supervisor/conf.d/00_set_system_timezone.conf
[program:set_system_timezone]
command = /scripts/set_system_timezone
startsecs = 0
autorestart = false
startretries = 1
